# pkgbuilds/bundle-xfce/files/etc/profile.d/xfce-core-meta-qt5ct.sh
[ ! -z "$XDG_SESSION_DESKTOP" ] && [ "$XDG_SESSION_DESKTOP" == "xfce" ] && [ -x /usr/bin/qt5ct ] && export QT_QPA_PLATFORMTHEME=qt5ct
