# pkgbuilds/dynu-ddns-updater/files/usr/lib/systemd/system/dynu-ddns-updater.service
[Unit]
Description=Update Dynu dynamic DNS

[Service]
Type=oneshot
Restart=no
RemainAfterExit=no
ExecStart=/usr/bin/dynu-ddns-updater
