# pkgbuilds/bundle-bluetooth/install.install
post_install() {
	## Enable bluetooth service
	systemctl --now enable bluetooth &> /dev/null
}
