# pkgbuilds/bundle-nvidia-pro-340xx/install.install
post_install() {
	systemctl enable nvidia-persistenced
}
